# modules/apache-maven/files/etc/profile.d/apache-maven.sh
export M3_HOME=/opt/apache-maven
export PATH=$PATH:$M3_HOME/bin
