# modules/mongodb/files/etc/profile.d/mongodb.sh
export MONGODB_HOME=/opt/mongodb
export PATH=$PATH:$MONGODB_HOME/bin
