# modules/redis/files/etc/profile.d/redis.sh
export REDIS_HOME=/opt/redis
export PATH=$PATH:$REDIS_HOME/bin
